# backend/systemd/cryptosignal-futures.service
[Unit]
Description=CryptoSignal Futures Worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cryptosignal-app/backend
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/cryptosignal-app/backend"
EnvironmentFile=/opt/cryptosignal-app/backend/.env
ExecStart=/opt/cryptosignal-app/backend/venv/bin/python3 workers/worker_futures.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptosignal-futures

[Install]
WantedBy=multi-user.target