# backend/systemd/cryptosignal-telegram-admin.service
[Unit]
Description=CryptoSignal Telegram Admin Bot
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cryptosignal-app/backend
Environment="PYTHONUNBUFFERED=1"
Environment="TELEGRAM_ADMIN_BOT_TOKEN=YOUR_ADMIN_BOT_TOKEN"
Environment="ADMIN_CHAT_ID=YOUR_CHAT_ID"
ExecStart=/usr/bin/python3 worker_telegram_admin.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptosignal-telegram-admin

[Install]
WantedBy=multi-user.target