# backend/systemd/cryptosignal-ai-analyst.service
[Unit]
Description=CryptoSignal AI Analyst Worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cryptosignal-app/backend
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 worker_ai_analyst.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptosignal-ai-analyst

[Install]
WantedBy=multi-user.target