# backend/systemd/cryptosignal-price-alerts.service
[Unit]
Description=CryptoSignal Price Alerts Worker
Documentation=https://github.com/onurcangnc/CryptoSignalApp
After=network.target redis.service cryptosignal-prices.service
Wants=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cryptosignal-app/backend
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/cryptosignal-app/backend"
EnvironmentFile=/opt/cryptosignal-app/backend/.env
ExecStart=/opt/cryptosignal-app/backend/venv/bin/python3 /opt/cryptosignal-app/backend/workers/worker_price_alerts.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptosignal-price-alerts

# Resource limits
MemoryMax=128M
CPUQuota=20%

[Install]
WantedBy=multi-user.target